# rgb steelseries keyboard power off at boot.service
[Unit]
Description=keyboard colour stop Service

[Service]
ExecStart=/usr/bin/msikey

[Install]
WantedBy=graphical.target
